# other_configurations/todo.conf
[program:todo]
command = /usr/local/bin/gunicorn -b 127.0.0.1:9000 run:app
directory = /Users/leo/Documents/WorkSpace/Python/Flask/Todo